# modularize/customized/croots/shared/etc/rc.local
#!/bin/bash

#add node starup setting here!
