# mod-stornode/customized/croots/stornode/etc/rc.local
#!/bin/bash
/etc/rc.d/init.d/syslog start
/etc/rc.d/init.d/sshd start

#mount bricks
mount -a

/etc/rc.d/init.d/glusterd start
